# systemd/honcho-embedding-endpoint-watchdog.timer
[Unit]
Description=Run Honcho embedding endpoint fallback watchdog every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true
Unit=honcho-embedding-endpoint-watchdog.service

[Install]
WantedBy=timers.target
